# roles/common/files/backup-local.service
[Unit]
Description=Perform backups locally.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-local
